# yggdrasill/hermes-laifu-monitor.service
[Unit]
Description=Hermes Lai.Fu Monitor

[Service]
Type=oneshot
ExecStart=/home/ken/hermes-mesh/yggdrasill/laifu-monitor.sh

[Install]
WantedBy=default.target
